# deploy/ww2survival.service
[Unit]
Description=WW2 Survival — Serveur dédié
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ww2survival
ExecStart=/opt/ww2survival/.venv/bin/python server_headless.py --autostart 15
Restart=on-failure
RestartSec=5

# SDL headless (pas d'écran physique)
Environment=SDL_VIDEODRIVER=dummy
Environment=SDL_AUDIODRIVER=dummy

# Sécurité : exécuter en tant qu'utilisateur non-root
User=ubuntu
Group=ubuntu

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
